# etc/profile.d/pdsh.csh
################################################################################
##
## Alces Underware - Shell configuration
## Copyright (c) 2008-2015 Alces Software Ltd
##
################################################################################
setenv GENDERS_FILE /var/lib/underware/rendered/system/genders
setenv PDSH_GENDERS_FILE $GENDERS_FILE
setenv PATH /opt/underware/opt/pdsh/bin:/opt/underware/opt/genders/bin:$PATH
